# init.cointhink
#!/sbin/openrc-run

name="cointhink script"

depend() {
	need net localmount
	after firewall
}

start() {
    ebegin "Starting cointhink script.py"
    start-stop-daemon --background --start --exec /usr/bin/python3 /cointhink/start.py \
    --make-pidfile --pidfile /cointhink/pid --env PYTHONPATH=/cointhink/proto 
    eend $?
}

